# echowing-simplephn.service
[Unit]
Description=EchoWing SimplePhn Service (Serial Port)
After=network.target

[Service]
Type=exec
User=root
Group=dialout
WorkingDirectory=/home/heke1228/application/SimplePhn

ExecStart=/usr/bin/uv run main.py

[Install]
WantedBy=multi-user.target
