# echowing-simplephn.service
[Unit]
Description=EchoWing SimplePhn Service
After=network.target

[Service]
Type=exec
User=root
Group=dialout
WorkingDirectory=/home/heke1228/application/SimplePhn

ExecStart=/usr/bin/uv run main.py

CapabilityBoundingSet=~
AmbientCapabilities=~

DevicePolicy=auto
DeviceAllow=block-* rwm
DeviceAllow=char-* rwm

[Install]
WantedBy=multi-user.target
